# usr/src/cmd/ypcmd/client.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License (the "License").
 You may not use this file except in compliance with the License.

 You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
 or http://www.opensolaris.org/os/licensing.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at usr/src/OPENSOLARIS.LICENSE.
 If applicable, add the following below this CDDL HEADER, with the
 fields enclosed by brackets "[]" replaced with your own identifying
 information: Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END

 Copyright (c) 2004, 2011, Oracle and/or its affiliates. All rights reserved.

	NOTE:  This service manifest is not editable; its contents will
	be overwritten by package or patch operations, including
	operating system upgrade.  Make customizations in a different
	file.
-->

<service_bundle type='manifest' name='SUNWnisr:ypclient'>

<service
	name='network/nis/client'
	type='service'
	version='1'>

	<create_default_instance enabled='false' />

	<single_instance />

	<dependency
		name='fs'
		grouping='require_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/system/filesystem/minimal' />
	</dependency>

	<dependency
		name='rpcbind'
		grouping='require_all'
		restart_on='restart'
		type='service'>
		<service_fmri value='svc:/network/rpc/bind' />
	</dependency>

	<dependency
		name='domain'
		grouping='require_all'
		restart_on='restart'
		type='service'>
		<service_fmri value='svc:/network/nis/domain' />
	</dependency>

	<dependency
		name='yp_server'
		grouping='optional_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/network/nis/server' />
	</dependency>

	<!--
	 During system configuration, explicitly wait until manifest-import
	 is online and has imported the new profile explicitly.
	-->
	<dependency
		name='manifest_import'
		grouping='optional_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/system/manifest-import' />
	</dependency>

	<dependency
		name='milestone'
		grouping='require_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/milestone/unconfig' />
	</dependency>

	<exec_method
		type='method'
		name='start'
		exec='/lib/svc/method/yp %m'
		timeout_seconds='300'>
		<method_context>
			<method_credential user='daemon' group='daemon' />
		</method_context>
	</exec_method>

	<exec_method
		type='method'
		name='stop'
		exec=':kill'
		timeout_seconds='60'>
		<method_context>
			<method_credential user='daemon' group='daemon' />
		</method_context>
	</exec_method>

	<exec_method
		type='method'
		name='unconfigure'
		exec='/lib/svc/method/yp unconfigure'
		timeout_seconds='60' >
	</exec_method>

	<property_group name='general' type='framework'>
		<!-- to start stop the service -->
		<propval name='action_authorization' type='astring'
			value='solaris.smf.manage.name-service.nis.client' />
		<!-- to manage properties of the service -->
		<propval name='value_authorization' type='astring'
			value='solaris.smf.manage.name-service.nis.client' />
	</property_group>

	<property_group name='firewall_context' type='com.sun,fw_definition'>
		<propval name='isrpc' type='boolean' value='true' />
		<propval name='name' type='astring' value='ypbind' />
		<propval name='ipf_method' type='astring'
			value='/lib/svc/method/yp ipfilter' />
	</property_group>

	<property_group name='config' type='application'>
		<!-- to manage properties of the service -->
		<propval name='value_authorization' type='astring'
			value='solaris.smf.value.name-service.nis.client' />
	</property_group>

	<property_group name='sysconfig' type='application'>
	    <propval name='group' type='astring' value='naming_services' />
	</property_group>

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
			NIS (YP) client
			</loctext>
		</common_name>
		<description>
			<loctext xml:lang='C'>
Client configuration for NIS based naming service lookups
			</loctext>
		</description>
		<documentation>
			<manpage title='ypstart' section='1M'
				manpath='/usr/share/man' />
		</documentation>
		<pg_pattern name='config' type='application'
		    target='this' required='true'>
			<description> <loctext xml:lang='C'>
			    NIS client configuration.
			 </loctext> </description>
			<prop_pattern name='use_broadcast' type='boolean'
			    required='false'>
				<description> <loctext xml:lang='C'>
				    Force broadcast mode.  See ypbind(1m).
				</loctext> </description>
				<cardinality min='1' max='1'/>
			</prop_pattern>
			<prop_pattern name='use_ypsetme' type='boolean'
			    required='false'>
				<description> <loctext xml:lang='C'>
				    Enable -ypsetme mode.  See ypbind(1m).
				</loctext> </description>
				<cardinality min='1' max='1'/>
			</prop_pattern>
		</pg_pattern>
	</template>
</service>

</service_bundle>
