# usr/src/cmd/shadowd/shadowd.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
 Copyright (c) 2011, Oracle and/or its affiliates. All rights reserved.

 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License (the "License").
 You may not use this file except in compliance with the License.

 You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
 or http://www.opensolaris.org/os/licensing.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at usr/src/OPENSOLARIS.LICENSE.
 If applicable, add the following below this CDDL HEADER, with the
 fields enclosed by brackets "[]" replaced with your own identifying
 information: Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END

	NOTE:  This service manifest is not editable; its contents will
	be overwritten by package or patch operations, including
	operating system upgrade.  Make customizations in a different
	file.

	NOTE: Do NOT edit the 'config' properties directly, use svccfg(1M)
	instead.

-->

<service_bundle type='manifest' name='fs-shadowd'>

<service
	name='system/filesystem/shadowd'
	type='service'
	version='1'>

	<create_default_instance enabled='false' />

	<single_instance />

	<!-- nfs is a best guess for now.  In theory we could  -->
	<!-- not need nfs to migrate if it's local to local.   -->
	<!-- Also in near future all migrations might come     -->
	<!-- via smb client and so whether nfs is up would be  -->
	<!-- uninteresting.                                    -->
	<dependency
	    name='fs'
	    type='service'
	    grouping='require_all'
	    restart_on='none'>
		<service_fmri value='svc:/system/filesystem/local' />
	</dependency>

	<exec_method
	    type='method'
	    name='start'
	    exec='/usr/lib/fs/shadowd'
	    timeout_seconds='0'>
	</exec_method>

	<exec_method
	    type='method'
	    name='refresh'
	    exec=':kill -HUP'
	    timeout_seconds='60'>
	</exec_method>

	<exec_method
	    type='method'
	    name='stop'
	    exec=':kill'
	    timeout_seconds='60'>
	</exec_method>

	<property_group name='config_params' type='application'>
	    <propval name='shadow_threads' type='integer' value='8' />
	    <propval name='shadow_throttle' type='integer' value='1000000' />
	</property_group>

	<stability value='Stable' />

	<template>
	    <common_name>
		<loctext xml:lang='C'>
		    Shadow Migration Control
		</loctext>
	    </common_name>
	    <documentation>
		<manpage title='shadowd' section='1M'
		    manpath='/usr/share/man' />
	    </documentation>
	</template>
</service>

</service_bundle>
